# iso-builder/overlays/et-v2-general/etc/skel/.config/systemd/user/et-supervisor.service
[Unit]
Description=LiaisonOS Process Supervisor
Documentation=https://liaisonos.com

[Service]
Type=simple
ExecStart=/opt/emcomm-tools/bin/et-supervisor
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
