# iso-builder/overlays/et-v2-general/etc/udev/rules.d/75-et-generic-remove.rules
# Author   : Gaston Gonzalez
# Date     : 14 October 2024
# Updated  : 24 October 2024
# Purpose  : udev rules for generic remove events
#
# Preconditions:
# 1. CAT control device (or DigiRig Lite) is unplugged
#
# Postconditions:
# 1. rigctld is stopped

SUBSYSTEM=="tty", ENV{ET_SUBDEVICE}=="CAT", \
  ACTION=="remove", \
  RUN+="/opt/emcomm-tools/sbin/wrapper-rigctld.sh stop"


# Note: Audio remove no longer stops rigctld — USB audio can briefly
# disconnect during Wine mode switches, causing a rigctld kill loop.
# Only CAT serial removal should stop rigctld.
