# iso-builder/overlays/et-v2-general/etc/skel/.config/autostart/et-auto-scale.desktop
[Desktop Entry]
Type=Application
Name=LiaisonOS Auto Scale
Name[fr]=LiaisonOS Échelle Auto
Comment=Automatically adjust display scaling based on screen DPI
Comment[fr]=Ajuste automatiquement l'échelle d'affichage selon le DPI de l'écran
Exec=/opt/emcomm-tools/bin/et-auto-scale-once
Icon=video-display
Terminal=false
Hidden=false
X-XFCE-Autostart-Override=true
StartupNotify=false
OnlyShowIn=XFCE;
