# iso-builder/overlays/et-v2-general/etc/skel/.config/autostart/onboard.desktop
[Desktop Entry]
Type=Application
Name=Onboard Virtual Keyboard
Exec=onboard --start-minimized
Hidden=false
X-GNOME-Autostart-enabled=true
