# iso-builder/overlays/et-v2-general/etc/skel/.config/autostart/xfce4-power-manager.desktop
[Desktop Entry]
Type=Application
Name=Power Manager
Exec=xfce4-power-manager
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
